# packaging/starter.path
[Unit]
Description=Path activation for the "starter" programme

[Path]
PathExists=/run/amd_ready
