# packaging/starter.path
[Unit]
Description=Start the "starter" program

[Path]
PathExists=/run/.wm_ready
